# s6/etc/services.d/pagermaid/run
#!/usr/bin/with-contenv bash

cd /pagermaid/workdir
exec s6-setuidgid pagermaid python -m pagermaid

